# docker-setup/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nextjs]
command = yarn start

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
numprocs=1